# roles/openqa/dispatcher/files/schedule-fcos-builds.cron
#!/bin/sh

# schedule openQA jobs for the current Fedora CoreOS builds in the
# various streams; if the current build has already been tested,
# this will not create any new jobs

for stream in stable testing next; do /usr/local/bin/fedora-openqa fcosbuild $stream > /dev/null 2>&1 || /usr/bin/true; done
