# roles/review-stats/build/files/review-stats-new.cron
0,30 * * * * apache /usr/local/bin/review-stats.py -d /srv/web/review-stats/
