# roles/review-stats/build/files/review-stats-new.cron
0,30 * * * * apache /usr/bin/review-stats -d /srv/web/review-stats/
