# roles/web-data-analysis/files/countme-update.cron
MAILTO=root@fedoraproject.org
0 09  * * * countme /usr/local/bin/countme-update.sh > /dev/null
