# roles/grobisplitter/files/cents8-split.cron
10 10  * * * root /usr/local/bin/centos-stream8-split.sh
